# hcloud-fip-sync.service
[Unit]
Description=hcloud floating-ip sync service

[Service]
Type=oneshot
ExecStart=/usr/bin/hcloud-fip-sync

[Install]
WantedBy=multi-user.target